# assets/setup/supervisor-postgresql.conf
[program:postgresql]
priority=20
directory=/tmp
command=/usr/lib/postgresql/9.3/bin/postgres -D /etc/postgresql/9.3/main
user=postgres
autostart=true
autorestart=true
stopsignal=INT
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
