# assets/setup/supervisor-postgresql.conf
[program:postgresql]
priority=20
directory=/tmp
command=/usr/lib/postgresql/9.3/bin/postgres -D /etc/postgresql/9.3/main
user=postgres
autostart=true
autorestart=true
stdout_logfile=${LOG_DIR}/supervisor/%(program_name)s.log
stderr_logfile=${LOG_DIR}/supervisor/%(program_name)s.log
